# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable daydash-service.service
systemctl start daydash-service.service
