# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop daydash-service.service
systemctl disable daydash-service.service
